# .github/server/guzhi-assistant.service
[Unit]
Description=估值助手 ASP.NET Core Service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/www/wwwroot/估值助手
ExecStart=/usr/bin/dotnet /www/wwwroot/估值助手/估值助手.dll --urls http://127.0.0.1:7084
Restart=always
RestartSec=3
KillSignal=SIGINT
SyslogIdentifier=guzhi_assistant
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
